# supervisord.conf
[supervisord]
nodaemon=true

[program:bluealsa-aplay]
directory=/usr/src/
command=/bin/bash start.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
